# v0.6.0-20260625_2022/.workbuddy/staging/root.cron
0 3 * * * /usr/local/bin/backup-oa.sh >> /var/log/oa-backup.log 2>&1
15 3 * * * /usr/local/bin/backup-oa-files.sh >> /var/log/oa-backup.log 2>&1
*/5 * * * * /usr/local/bin/oa-monitor.sh >> /var/log/oa-monitor.log 2>&1
